# Dockerfiles/cowire/cowrie.conf
[program:cowrie]
command=/home/cowrie/cowrie/bin/cowrie start
autostart=true
autorestart=true
user=cowrie
stdout_logfile=/var/log/cowrie/cowrie.log
stderr_logfile=/var/log/cowrie/cowrie.err