# Dockerfiles/cowire/supervisord.conf
[supervisord]
nodaemon=true

[include]
files = /etc/supervisor/conf.d/*.conf

[program:cowrie]
command=cowrie start